# init.d/gifbox.sh
#! /bin/sh
### BEGIN INIT INFO
# Provides:          gifbox
# Required-Start:    $all
# Required-Stop:     
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Manage gifbox
### END INIT INFO

/home/pi/gifbox/gifbox.sh $1